# installer/targets/megatools/mega.timer
[Unit]
Description=Mega Sync Timer (runs every minute)
[Timer]
OnUnitActiveSec=1m
Persistent=true
Unit=mega.service
[Install]
WantedBy=timers.target
